# documize.service
[Unit]
Description=Documize the free Confluence Alternative
Requires=network.target
After=postrgresql.service

[Service]
Type=simple
User=vagrant
WorkingDirectory=/opt/documize/
ExecStart=/opt/documize/documize-community-linux-amd64 settings.toml

[Install]
WantedBy=multi-user.target
