# authelia/authelia.service
[Unit]
Description=Authelia authentication and authorization server
After=multi-user.target

[Service]
User=authelia
Group=authelia
ExecStart=/usr/local/bin/authelia --config /etc/authelia/config.yml
SyslogIdentifier=authelia

[Install]
WantedBy=multi-user.target
